# ansible/GrafanaPrometheus/monitor/files/node_exporter/node_exporter.service
[Unit]
Description=Node Exporter

[Service]
User=root
ExecStart=/usr/sbin/node_exporter

[Install]
WantedBy=multi-user.target